# deployment/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# برنامه اصلی مدیریت ربات‌ها - Main Bot Management Program
[program:telegram-bots-main]
command=/app/venv/bin/python unified_bot_launcher.py
directory=/app
user=telegrambot
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stderr_logfile=/app/logs/main_error.log
stdout_logfile=/app/logs/main_output.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
stderr_logfile_backups=5
stdout_logfile_backups=5
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app
priority=100

# ربات مانیتورینگ - Monitoring Bot
[program:telegram-bots-monitor]
command=/app/venv/bin/python monitoring_bot.py
directory=/app
user=telegrambot
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=10
stderr_logfile=/app/logs/monitor_error.log
stdout_logfile=/app/logs/monitor_output.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app
priority=200

# ربات گزارش‌دهی - Report Bot
[program:telegram-bots-report]
command=/app/venv/bin/python report_bot.py
directory=/app
user=telegrambot
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=10
stderr_logfile=/app/logs/report_error.log
stdout_logfile=/app/logs/report_output.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app
priority=300

# سرویس وب پنل - Web Panel Service
[program:web-panel]
command=/app/venv/bin/python start_system.py
directory=/app
user=telegrambot
autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=15
stderr_logfile=/app/logs/web_error.log
stdout_logfile=/app/logs/web_output.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app
priority=400

# سرویس مانیتورینگ سیستم - System Monitoring Service
[program:system-monitor]
command=/app/venv/bin/python deployment/monitoring.py --interval 30
directory=/app
user=telegrambot
autostart=true
autorestart=true
startretries=5
startsecs=15
stopwaitsecs=10
stderr_logfile=/app/logs/system_monitor_error.log
stdout_logfile=/app/logs/system_monitor_output.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app
priority=500

# سرویس بک‌آپ - Backup Service
[program:backup-service]
command=/app/venv/bin/python deployment/backup_service.py
directory=/app
user=telegrambot
autostart=false
autorestart=false
stderr_logfile=/app/logs/backup_error.log
stdout_logfile=/app/logs/backup_output.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile_backups=5
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app
priority=600

# گروه‌بندی برنامه‌ها - Program Groups
[group:core-services]
programs=telegram-bots-main,telegram-bots-monitor,telegram-bots-report
priority=100

[group:web-services]
programs=web-panel
priority=200

[group:support-services]
programs=system-monitor
priority=300

[group:maintenance-services]
programs=backup-service
priority=400

# تنظیمات Event Listener برای مدیریت رویدادها
[eventlistener:crash-handler]
command=/app/deployment/scripts/crash_handler.py
events=PROCESS_STATE_FATAL,PROCESS_STATE_EXITED
stderr_logfile=/app/logs/crash_handler_error.log
stdout_logfile=/app/logs/crash_handler_output.log
environment=PYTHONUNBUFFERED=1,PYTHONPATH=/app

# تنظیمات امنیتی و محدودیت‌ها
[include]
files = /app/deployment/supervisor/conf.d/*.conf